# terraform/docker-nextcloud.cron
#!/bin/sh
/etc/rc.d/init.d/docker-nextcloud update
