# terraform/files/1post-yum-security.cron
#!/bin/bash

if /sbin/grubby --default-kernel | /bin/grep -q "$(/bin/uname -r)"
then
  /etc/rc.d/init.d/docker-nextcloud start
else
  /sbin/reboot
fi

