# terraform/files/1post-yum-security.cron
#!/bin/bash
/etc/rc.d/init.d/docker-nextcloud start
/sbin/grubby --default-kernel | /bin/grep -q "$(/bin/uname -r)" || /sbin/reboot
