# terraform/files/1post-yum-security.cron
#!/bin/bash

is_running() {
  test -n "$(/etc/init.d/docker-nextcloud status)"
}

if /usr/bin/needs-restarting -r
then
  is_running || /etc/rc.d/init.d/docker-nextcloud restart
else
  /sbin/reboot
fi

