# terraform/files/1post-yum-security.cron
#!/bin/bash

if /usr/bin/needs-restarting -r
then
  /etc/rc.d/init.d/docker-nextcloud start
else
  /sbin/reboot
fi

